# deploy/pdf-editor.service
[Unit]
Description=PDF Editor (RunPro) web service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/pdf-editor
# Учётки и host/port берутся из .env (PDF_EDITOR_HOST=127.0.0.1, PDF_EDITOR_PORT=8765)
EnvironmentFile=/opt/pdf-editor/.env
ExecStart=/opt/pdf-editor/venv/bin/python pdf_editor.py gui --no-browser
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
